# systemd/neohome-minutely.sample.timer
[Unit]
Description=Timer for Neohome minutely script

[Timer]
OnCalendar=minutely
AccuracySec=1us

[Install]
WantedBy=timers.target
