# systemd/phpmd-minutely.sample.timer
[Unit]
Description=Timer for PHPMD minutely script

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
