# service/systemd/shcc.target
[Unit]
Wants=shcc-minutely.timer

[Install]
WantedBy=default.target
