# service/systemd/shcc.target
[Unit]
Wants=shcc-minutely.timer shcc@xiaomi.service shcc@yeelight.service shcc@miio.service shcc@tts.service

[Install]
WantedBy=default.target
